# Syntropic System/System Infrastructure/LaunchAgents/com.obsidian.dailybrief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.obsidian.dailybrief</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/elbowspeak1/Documents/Obsidian Vault/Doikayt/Syntropic System/1. Executive & Portfolio Management/Automation/daily-brief-runner.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday -->
        <dict>
            <key>Hour</key>
            <integer>5</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>1</integer>
        </dict>
        <!-- Tuesday -->
        <dict>
            <key>Hour</key>
            <integer>5</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>2</integer>
        </dict>
        <!-- Wednesday -->
        <dict>
            <key>Hour</key>
            <integer>5</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>3</integer>
        </dict>
        <!-- Thursday -->
        <dict>
            <key>Hour</key>
            <integer>5</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>4</integer>
        </dict>
        <!-- Friday -->
        <dict>
            <key>Hour</key>
            <integer>5</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>5</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/elbowspeak1/Documents/Obsidian Vault/Doikayt/.cache/logs/launchagent-out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/elbowspeak1/Documents/Obsidian Vault/Doikayt/.cache/logs/launchagent-error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
